# supervisord.conf
[supervisord]
nodaemon=true

[program:xrdp]
command=/usr/sbin/xrdp --nodaemon
autostart=true
autorestart=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
